# salt/duplicity/duplicity.service
[Unit]
Description=Backup

[Service]
Nice=19
IOSchedulingClass=2
IOSchedulingPriority=7
ExecStart=/root/bin/backup.sh
